# lib/systemd/system/ips-grouping.service
[Unit]
Description=ips-grouping

[Service]
User=ipstudio
ExecStart=/usr/bin/ips-grouping
